# target-root/etc/sharpinit/units/getty.target
[Unit]
Description=getty.target

[Install]
WantedBy=default.target
